# php-fpm.conf
description     "Start php-fpm worker pool for a Fast-CGI server"

# Start when networking starts
start on started networking
stop on shutdown

exec /usr/local/php-5.4.5/sbin/php-fpm
